# systemd/tilda.service
[Unit]
Description=Tilda

[Service]
ExecStart=/usr/bin/tilda

[Install]
WantedBy=multi-user.target
